# roles/videoserver/templates/rtsp-simple-server.service.j2
[Unit]
After=network.target

[Service]
ExecStart={{ videoserver_executable_location }} {{ videoserver_config_location }}

[Install]
WantedBy=multi-user.target
