# roles/freetakserver/templates/fts.service.j2
[Unit]
Description=FreeTAKServer service
After=network.target
StartLimitIntervalSec=0

[Service]
Type=simple
Restart=always
RestartSec=1
User=fts
Group=fts
StandardOutput=append:/var/log/{{ fts_service_name }}/{{ fts_service_name }}-stdout.log
StandardError=append:/var/log/{{ fts_service_name }}/{{ fts_service_name }}-stderr.log
Environment="FTS_FIRST_START=false"
ExecStart={{ fts_venv }}/bin/python3 -m FreeTAKServer.controllers.services.FTS

[Install]
WantedBy=multi-user.target
